# etc/udev/rules.d/90-hdmi.rules
SUBSYSTEM=="drm", ACTION=="change", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/run/user/1000/Xauthority", RUN+="/home/$USER/.local/bin/scripts/montoggle"
